# Curso-Introducción a Terminal y Línea de Comandos/crontab
minuto 0-59
hora 0-23
dia mes 1-31
mes 1-12
dia semana: 0-7 (0 y 7 domingo)
script/comando

ejemplo 1 col (min)
1
1,10,18
*/5
1-10
*

